# docker/ros.target
[Unit]
Description=ros
Requires=ros.target novncd@recipelab.service
After=ros.target novncd@recipelab.service
AllowIsolate=yes